# open-weather-fetch.service
[Unit]
Description="Weather data fetching"

[Service]
Type=simple
ExecStart=@bindir@/weather-fetch
Restart=on-failure

[Install]
WantedBy=default.target
